# flower/etc/service/flower/run
#!/bin/sh -e
exec flower --broker=$CELERY_BROKER_URL 2>&1
